# .config/systemd/user/check_fwupd_status.service
[Unit]
Description=Check whether firmware updates are available

[Service]
Type=oneshot
ExecStart=%h/.local/bin/check_fwupd_status.sh
