# .config/systemd/user/check_fwupd_status.timer
[Unit]
Description=Check whether firmware updates are available

[Timer]
OnStartupSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
